# scripts/ovirt-firstboot
#!/bin/bash
#
# ovirt-firstboot Puts the node into configuration mode.
#
# chkconfig: 2345 99 01
# description: ovirt-firstboot node configuration script
#
# Copyright (C) 2008 Red Hat, Inc.
# Written by Darryl L. Pierce <dpierce@redhat.com>
#
# This program is free software; you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published by
# the Free Software Foundation; version 2 of the License.
#
# This program is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU General Public License for more details.
#
# You should have received a copy of the GNU General Public License
# along with this program; if not, write to the Free Software
# Foundation, Inc., 51 Franklin Street, Fifth Floor, Boston,
# MA  02110-1301, USA.  A copy of the GNU General Public License is
# also available at http://www.gnu.org/copyleft/gpl.html.

## Source functions library
. /etc/init.d/functions
. /etc/init.d/ovirt-functions

start ()
{
    plymouth --hide-splash

    if is_auto_install; then
        ovirt-config-networking AUTO
        ovirt-config-storage AUTO
        ovirt-config-logging AUTO
    else
        ovirt-config-setup
    fi

    plymouth --show-splash
}

case "$1" in
    start)
        if is_managed; then
            exit 0
        fi

        printf "Starting ovirt-firstboot: "

        {
            start
        }

        test $? == 0 && success || failure
        echo
        ;;

    *)
        echo "Usage: ovirt-firstboot {start}"
        exit 2
esac

