# scripts/ovirt-firstboot
#!/bin/bash
#
# ovirt-firstboot Puts the node into configuration mode.
#
# chkconfig: 2345 99 01
# description: ovirt-firstboot node configuration script
#
# Copyright (C) 2008 Red Hat, Inc.
# Written by Darryl L. Pierce <dpierce@redhat.com>
#
# This program is free software; you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published by
# the Free Software Foundation; version 2 of the License.
#
# This program is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU General Public License for more details.
#
# You should have received a copy of the GNU General Public License
# along with this program; if not, write to the Free Software
# Foundation, Inc., 51 Franklin Street, Fifth Floor, Boston,
# MA  02110-1301, USA.  A copy of the GNU General Public License is
# also available at http://www.gnu.org/copyleft/gpl.html.

## Source functions library
. /etc/init.d/functions
. /etc/init.d/ovirt-functions

prog=ovirt-firstboot
lockfile=/var/lock/subsys/$prog

trap '__st=$?; stop_log; exit $__st' 0
trap 'exit $?' 1 2 13 15

start ()
{
    if ! is_firstboot && ! is_auto_install && ! is_upgrade; then
        return
    fi
	
	touch $lockfile

    # by default in standalone, disable collectd until server is configured
    service collectd stop > /dev/null 2>&1
    chkconfig collectd off

    if is_auto_install; then
        ovirt-config-storage AUTO
        ovirt-config-hostname AUTO
        ovirt-config-networking AUTO
        ovirt-config-logging AUTO
        ovirt-config-collectd AUTO
        ovirt-config-password AUTO
    fi

    if is_upgrade; then
        mount_live
        ovirt-config-boot /live "$OVIRT_BOOTPARAMS" no
        disable_firstboot
        reboot
        return
    fi

    if is_firstboot; then
        plymouth --hide-splash

        ovirt-config-setup -x < /dev/console

        plymouth --show-splash
    fi
    disable_firstboot

    ovirt_store_firstboot_config >> $OVIRT_LOGFILE 2>&1

    rm -f $lockfile
}

case "$1" in
    start)
        if is_managed; then
            exit 0
        fi

        printf "Starting ovirt-firstboot: "

        start

        test $? == 0 && success || failure
        echo
        ;;
    status)
        status $prog
        ;;
    reload)
        stop
        start
        ;;
    stop)
        stop
        ;;
    *)
        echo "Usage: ovirt-firstboot {start}"
        exit 2
esac

